# deploy/ansible/roles/maxwell_daemon/templates/maxwell-daemon.service.j2
[Unit]
Description=Maxwell-Daemon agent orchestrator
After=network.target
Wants=network.target

[Service]
Type=simple
User={{ maxwell_user }}
Group={{ maxwell_group }}
WorkingDirectory={{ maxwell_data_dir }}
ExecStart={{ maxwell_venv }}/bin/maxwell-daemon serve \
    --config {{ maxwell_config_dir }}/config.toml
Restart=on-failure
RestartSec=5s

# Harden the service
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths={{ maxwell_data_dir }} {{ maxwell_log_dir }}

[Install]
WantedBy=multi-user.target
